# misc/mdb.supervisor.ini
[program:mdb]
process_name = mdb
command = /sites/mdb/mdb server
directory = /sites/mdb/
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/mdb/logs/mdb.log

[program:mdb_gor]
process_name = mdb_gor
command = /sites/mdb/mdb_gor.sh
directory = /sites/mdb/
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/mdb/logs/requests.log
